# pi/systemd/ghostroll-watch.service
[Unit]
Description=GhostRoll watch (SD ingest pipeline)
After=ghostroll-wifi-setup.service network-online.target ghostroll-firstboot.service
Wants=ghostroll-wifi-setup.service network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
EnvironmentFile=-/etc/ghostroll.env
# Do NOT use bash-style default expansion (":-") here; systemd does not evaluate it.
# GhostRoll already reads GHOSTROLL_SD_LABEL with a default of "auto-import".
ExecStart=/usr/local/bin/ghostroll watch
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target


