# pi/systemd/ghostroll-update.service
[Unit]
Description=GhostRoll auto-update from Git remote
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/ghostroll.env
ExecStart=/usr/local/sbin/ghostroll-update.sh




